# deployment/keycloak/keycloak.service
[Unit]
Description=Recipes KeyCloak service

[Service]
ExecStart=docker run keycloak
User=pi

[Install]
WantedBy=multi-user.target
